# cookbooks/god/templates/default/god-init.erb
#!/bin/sh


### BEGIN INIT INFO
# Provides:             god
# Required-Start:       $all
# Required-Stop:        $all
# Default-Start:        2 3 4 5
# Default-Stop:         0 1 6
# Short-Description:    God
### END INIT INFO


NAME=god
DESC=god


set -e


# Make sure the binary and the config file are present before proceeding
test -x /opt/ruby-enterprise/bin/god || exit 0


# Create this file and put in a variable called GOD_CONFIG, pointing to
# your God configuration file
GOD_CONFIG=/etc/god/master.god

. /lib/lsb/init-functions


RETVAL=0


case "$1" in
  start)
    echo -n "Starting $DESC: "
    /opt/ruby-enterprise/bin/god -c $GOD_CONFIG -P /var/run/god.pid -l /var/log/god.log
    RETVAL=$?
    echo "$NAME."
    ;;
  stop)
    echo -n "Stopping $DESC: "
    kill `cat /var/run/god.pid`
    RETVAL=$?
    echo "$NAME."
    ;;
  restart)
    echo -n "Restarting $DESC: "
    kill `cat /var/run/god.pid`
    /opt/ruby-enterprise/bin/god -c $GOD_CONFIG -P /var/run/god.pid -l /var/log/god.log
    RETVAL=$?
    echo "$NAME."
    ;;
  status)
    /opt/ruby-enterprise/bin/god status
    RETVAL=$?
    ;;
  *)
    echo "Usage: god {start|stop|restart|status}"
    exit 1
    ;;
esac

exit $RETVAL

